# deploy/server/systemd/ogrodio-cf-ufw.service
[Unit]
Description=Ogrodio — odśwież reguły UFW dla IP Cloudflare

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/ogrodio/app/deploy/ufw-cloudflare-only.sh
StandardOutput=append:/var/log/ogrodio/cf-ufw.log
StandardError=append:/var/log/ogrodio/cf-ufw.log

[Install]
WantedBy=multi-user.target
